# 3letnik/KP/koncna-projektna/etcService.service
[Unit]
Description=Run etcd conf

[Service]
ExecStart=/home/sk01/etconf.sh

[Install]
WantedBy=multi-user.target